# var/lib/dpkg/info/live-config-sysvinit.list
/.
/etc
/etc/init.d
/etc/init.d/live-config
/lib
/lib/live
/lib/live/config
/lib/live/config/0160-sysvinit
/lib/live/config/0180-sysv-rc
/usr
/usr/share
/usr/share/doc
/usr/share/doc/live-config-sysvinit
/usr/share/doc/live-config-sysvinit/changelog.gz
/usr/share/doc/live-config-sysvinit/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/live-config-sysvinit
